# deploy/claudeforalice.service
[Unit]
Description=Claude for Alice AI skill (Flask/gunicorn)
After=network.target

[Service]
Type=simple
User=claudeforalice
Group=claudeforalice
WorkingDirectory=/srv/claudeforalice

# Load environment variables from .env
EnvironmentFile=/srv/claudeforalice/.env

# Use gunicorn from the venv
ExecStart=/srv/claudeforalice/.venv/bin/gunicorn \
    app:app \
    --bind 127.0.0.1:${FLASK_PORT} \
    --workers 1 \
    --threads 8 \
    --worker-class gthread \
    --timeout 60 \
    --access-logfile /var/log/claudeforalice/access.log \
    --error-logfile  /var/log/claudeforalice/error.log \
    --log-level info

Restart=always
RestartSec=3

# Harden the service
# Note: NoNewPrivileges is intentionally omitted so the dashboard can run
# "sudo systemctl restart claudeforalice" via the web UI.
ProtectSystem=strict
ReadWritePaths=/var/log/claudeforalice /srv/claudeforalice /run/sudo
PrivateTmp=true

[Install]
WantedBy=multi-user.target
